# rtirq-resume.service
[Unit]
Description=Realtime IRQ thread system tuning (restart after suspend)
After=multi-user.target suspend.target

[Service]
Type=simple
ExecStart=/usr/sbin/rtirq start

[Install]
WantedBy=multi-user.target suspend.target
